# scripts/prod/backup/backup.timer
[Unit]
Description=timer exec backup.sh

[Timer]
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
